# pkg/deb/prerm
#!/bin/sh

update-rc.d -f kauth remove >/dev/null 2>&1 || : 

# uinstall
if [ -L /etc/init.d/kauth ]; then
    /etc/init.d/kauth stop >/dev/null 2>&1 || :
    rm -f /etc/init.d/kauth || :
fi
